# files/etc/uci-defaults/99-beryl7-defaults
#!/bin/sh
set -eu

uci -q batch <<'EOF'
set system.@system[0].hostname='Beryl7'
set system.@system[0].zonename='Asia/Shanghai'
set system.@system[0].timezone='CST-8'
commit system
EOF

if ! uci -q get network.lan >/dev/null; then
	uci -q set network.lan='interface'
fi

uci -q batch <<'EOF'
set network.lan.proto='static'
set network.lan.ipaddr='192.168.18.1'
set network.lan.netmask='255.255.255.0'
commit network
EOF

if ! uci -q get dhcp.@dnsmasq[0] >/dev/null; then
	uci -q add dhcp dnsmasq >/dev/null
fi

uci -q batch <<'EOF'
set dhcp.@dnsmasq[0].domainneeded='1'
set dhcp.@dnsmasq[0].boguspriv='1'
set dhcp.@dnsmasq[0].filterwin2k='0'
set dhcp.@dnsmasq[0].localise_queries='1'
set dhcp.@dnsmasq[0].rebind_protection='1'
set dhcp.@dnsmasq[0].rebind_localhost='1'
set dhcp.@dnsmasq[0].local='/lan/'
set dhcp.@dnsmasq[0].domain='lan'
set dhcp.@dnsmasq[0].expandhosts='1'
set dhcp.@dnsmasq[0].authoritative='1'
set dhcp.@dnsmasq[0].readethers='1'
set dhcp.@dnsmasq[0].leasefile='/tmp/dhcp.leases'
set dhcp.@dnsmasq[0].resolvfile='/tmp/resolv.conf.d/resolv.conf.auto'
set dhcp.@dnsmasq[0].localservice='1'
set dhcp.lan='dhcp'
set dhcp.lan.interface='lan'
set dhcp.lan.start='100'
set dhcp.lan.limit='150'
set dhcp.lan.leasetime='12h'
set dhcp.lan.dhcpv4='server'
set dhcp.lan.ignore='0'
commit dhcp
EOF

if [ -x /etc/init.d/dnsmasq ]; then
	/etc/init.d/dnsmasq enable
fi

if uci -q get luci.main >/dev/null; then
	uci -q batch <<'EOF'
set luci.main.lang='zh-cn'
commit luci
EOF
fi

if uci -q get firewall.@defaults[0] >/dev/null; then
	uci -q batch <<'EOF'
set firewall.@defaults[0].flow_offloading='1'
set firewall.@defaults[0].flow_offloading_hw='1'
commit firewall
EOF
fi

# 中文：首次启动时启用默认无线，并使用用户指定的 SSID 与密码。
setup_default_wifi() {
	local wifi_key_file wifi_key wifi_ssid

	[ -s /etc/config/wireless ] || return 0
	uci -q get wireless.default_radio0 >/dev/null || return 0

	wifi_key_file='/etc/beryl7-wifi-password'
	wifi_key='77915558'
	wifi_ssid='MT3600BE'
	umask 077
	printf '%s\n' "${wifi_key}" > "${wifi_key_file}"

	if uci -q get wireless.radio0 >/dev/null; then
		uci -q set wireless.radio0.country='CN'
	fi

	uci -q set wireless.default_radio0.disabled='0'
	uci -q set wireless.default_radio0.ssid="${wifi_ssid}"
	uci -q set wireless.default_radio0.encryption='psk2'
	uci -q set wireless.default_radio0.key="${wifi_key}"

	if uci -q get wireless.radio1 >/dev/null && uci -q get wireless.default_radio1 >/dev/null; then
		uci -q set wireless.radio1.country='CN'
		uci -q set wireless.default_radio1.disabled='0'
		uci -q set wireless.default_radio1.ssid="${wifi_ssid}"
		uci -q set wireless.default_radio1.encryption='psk2'
		uci -q set wireless.default_radio1.key="${wifi_key}"
	fi

	uci -q commit wireless
	# 中文：首次启动脚本运行较早，只写入 UCI 配置；系统后续启动 network/wifi 时会自动应用。
	# 中文：避免在这里主动 reload 无线，防止驱动初始化阶段卡住首次启动流程。
}

setup_default_wifi

# 中文：AdGuardHome 不接管 53 端口，避免影响 dnsmasq 的 DNS 与 DHCP。
# 中文：如用户已经初始化过 AdGuardHome 配置，只把 DNS 监听端口固定到 5335。
setup_adguardhome_dns_port() {
	local agh_config

	agh_config="$(uci -q get adguardhome.config.config_file 2>/dev/null || true)"
	[ -n "${agh_config}" ] || agh_config='/etc/adguardhome/adguardhome.yaml'

	if uci -q get adguardhome.config >/dev/null; then
		uci -q set adguardhome.config.enabled='0'
		uci -q set adguardhome.config.config_file="${agh_config}"
		uci -q commit adguardhome
	fi

	[ -s "${agh_config}" ] || return 0

	awk '
		BEGIN { in_dns = 0; done = 0 }
		/^dns:[[:space:]]*$/ { in_dns = 1; print; next }
		in_dns && /^[^[:space:]][^:]*:/ {
			if (!done) {
				print "  port: 5335"
				done = 1
			}
			in_dns = 0
		}
		in_dns && /^[[:space:]]+port:[[:space:]]*/ {
			if (!done) {
				print "  port: 5335"
				done = 1
			}
			next
		}
		{ print }
		END {
			if (in_dns && !done) {
				print "  port: 5335"
			}
		}
	' "${agh_config}" > "${agh_config}.tmp" && mv "${agh_config}.tmp" "${agh_config}"
}

setup_adguardhome_dns_port

patch_luci_cpuusage_backend() {
	local luci_ucode tmp_file

	luci_ucode='/usr/share/rpcd/ucode/luci'
	[ -s "${luci_ucode}" ] || return 0
	[ -x /sbin/cpuusage ] || return 0
	grep -q '/sbin/cpuusage' "${luci_ucode}" && return 0

	tmp_file="${luci_ucode}.tmp"
	# 中文：OpenWrt 主线默认只返回 CPU 百分比；这里复用 /sbin/cpuusage 输出 CPU 与 PPE/NPU 流表统计。
	awk '
		/const fd = popen\('\''top -n1 \| awk/ {
			print "\t\t\tconst fd = popen('\''/sbin/cpuusage'\'');"
			next
		}
		{ print }
	' "${luci_ucode}" > "${tmp_file}" && mv "${tmp_file}" "${luci_ucode}"
}

patch_luci_cpuusage_backend

if [ -x /etc/init.d/irqbalance ]; then
	/etc/init.d/irqbalance enable
fi

exit 0
