# files/etc/uci-defaults/99-beryl7-defaults
#!/bin/sh
set -eu

uci -q batch <<'EOF'
set system.@system[0].hostname='Beryl7'
set system.@system[0].zonename='Asia/Shanghai'
set system.@system[0].timezone='CST-8'
commit system
EOF

setup_build_version_description() {
	# 中文：构建时间由 CI 写入 /etc/beryl7-build-info；首次启动时追加到 LuCI 固件版本显示。
	local release_file version_file os_release_file build_label current_desc new_desc current_version current_os_release new_os_release

	[ -r /etc/beryl7-build-info ] || return 0
	. /etc/beryl7-build-info
	build_label="${BUILD_LABEL:-}"
	[ -n "${build_label}" ] || return 0

	release_file='/etc/openwrt_release'
	version_file='/etc/openwrt_version'
	os_release_file='/usr/lib/os-release'

	if [ -r "${release_file}" ]; then
		current_desc="$(
			. "${release_file}" 2>/dev/null
			printf '%s' "${DISTRIB_DESCRIPTION:-}"
		)"
		if [ -n "${current_desc}" ]; then
			case "${current_desc}" in
			*"${build_label}"*) ;;
			*)
				new_desc="${current_desc} / ${build_label}"
				sed -i "s|^DISTRIB_DESCRIPTION=.*|DISTRIB_DESCRIPTION='${new_desc}'|" "${release_file}"
				;;
			esac
		fi
	fi

	if [ -r "${version_file}" ]; then
		current_version="$(cat "${version_file}" 2>/dev/null)"
		case "${current_version}" in
		*"${build_label}"*) ;;
		*) printf '%s / %s\n' "${current_version}" "${build_label}" > "${version_file}" ;;
		esac
	fi

	if [ -r "${os_release_file}" ]; then
		current_os_release="$(
			awk -F= '$1 == "OPENWRT_RELEASE" {
				value = $2
				sub(/^"/, "", value)
				sub(/"$/, "", value)
				print value
			}' "${os_release_file}"
		)"
		if [ -n "${current_os_release}" ]; then
			case "${current_os_release}" in
			*"${build_label}"*) ;;
			*)
				new_os_release="${current_os_release} / ${build_label}"
				awk -v release="${new_os_release}" '
					/^OPENWRT_RELEASE=/ {
						print "OPENWRT_RELEASE=\"" release "\""
						next
					}
					{ print }
				' "${os_release_file}" > /tmp/os-release.beryl7
				mv /tmp/os-release.beryl7 "${os_release_file}"
				;;
			esac
		fi
	fi
}

setup_build_version_description

if ! uci -q get network.lan >/dev/null; then
	uci -q set network.lan='interface'
fi

uci -q batch <<'EOF'
set network.lan.proto='static'
set network.lan.ipaddr='192.168.16.1'
set network.lan.netmask='255.255.255.0'
commit network
EOF

if ! uci -q get dhcp.@dnsmasq[0] >/dev/null; then
	uci -q add dhcp dnsmasq >/dev/null
fi

uci -q batch <<'EOF'
set dhcp.@dnsmasq[0].domainneeded='1'
set dhcp.@dnsmasq[0].boguspriv='1'
set dhcp.@dnsmasq[0].filterwin2k='0'
set dhcp.@dnsmasq[0].localise_queries='1'
set dhcp.@dnsmasq[0].rebind_protection='1'
set dhcp.@dnsmasq[0].rebind_localhost='1'
set dhcp.@dnsmasq[0].local='/lan/'
set dhcp.@dnsmasq[0].domain='lan'
set dhcp.@dnsmasq[0].expandhosts='1'
set dhcp.@dnsmasq[0].authoritative='1'
set dhcp.@dnsmasq[0].readethers='1'
set dhcp.@dnsmasq[0].leasefile='/tmp/dhcp.leases'
set dhcp.@dnsmasq[0].resolvfile='/tmp/resolv.conf.d/resolv.conf.auto'
set dhcp.@dnsmasq[0].localservice='1'
set dhcp.lan='dhcp'
set dhcp.lan.interface='lan'
set dhcp.lan.start='100'
set dhcp.lan.limit='150'
set dhcp.lan.leasetime='12h'
set dhcp.lan.dhcpv4='server'
set dhcp.lan.ignore='0'
commit dhcp
EOF

if [ -x /etc/init.d/dnsmasq ]; then
	/etc/init.d/dnsmasq enable
fi

setup_zerotier_no_dhcp() {
	# 中文：ZeroTier 接口地址由 zerotier-one 托管，不在虚拟网卡上启动 DHCP。
	if uci -q get network.zerotier >/dev/null; then
		uci -q set network.zerotier.proto='none'
		uci -q commit network
	fi

	if uci -q get dhcp >/dev/null; then
		uci -q set dhcp.zerotier='dhcp'
		uci -q set dhcp.zerotier.interface='zerotier'
		uci -q set dhcp.zerotier.ignore='1'
		uci -q commit dhcp
	fi
}

setup_zerotier_no_dhcp

if uci -q get luci.main >/dev/null; then
	uci -q batch <<'EOF'
set luci.main.lang='zh-cn'
commit luci
EOF
fi

if uci -q get firewall.@defaults[0] >/dev/null; then
	uci -q batch <<'EOF'
set firewall.@defaults[0].flow_offloading='1'
set firewall.@defaults[0].flow_offloading_hw='1'
commit firewall
EOF
fi

setup_default_wifi() {
	# 中文：首次启动时启用默认无线，并使用用户指定的 SSID 与密码。
	local wifi_key_file wifi_key wifi_ssid_2g wifi_ssid_5g

	[ -s /etc/config/wireless ] || return 0
	uci -q get wireless.default_radio0 >/dev/null || return 0

	wifi_key_file='/etc/beryl7-wifi-password'
	wifi_key='77915558'
	wifi_ssid_2g='MT3600BE_IOT'
	wifi_ssid_5g='MT3600BE'
	umask 077
	printf '%s\n' "${wifi_key}" > "${wifi_key_file}"

	if uci -q get wireless.radio0 >/dev/null; then
		uci -q set wireless.radio0.country='US'
		uci -q set wireless.radio0.channel='11'
		uci -q set wireless.radio0.htmode='HT20'
		uci -q set wireless.radio0.txpower='23'
		uci -q set wireless.radio0.cell_density='0'
		uci -q set wireless.radio0.legacy_rates='1'
	fi

	uci -q set wireless.default_radio0.disabled='0'
	uci -q set wireless.default_radio0.ssid="${wifi_ssid_2g}"
	uci -q set wireless.default_radio0.encryption='psk2+ccmp'
	uci -q set wireless.default_radio0.key="${wifi_key}"
	uci -q delete wireless.default_radio0.ieee80211w 2>/dev/null || true

	if uci -q get wireless.radio1 >/dev/null && uci -q get wireless.default_radio1 >/dev/null; then
		uci -q set wireless.radio1.country='US'
		uci -q set wireless.radio1.channel='149'
		uci -q set wireless.radio1.htmode='EHT80'
		uci -q set wireless.radio1.txpower='30'
		uci -q set wireless.radio1.cell_density='0'
		uci -q set wireless.default_radio1.disabled='0'
		uci -q set wireless.default_radio1.ssid="${wifi_ssid_5g}"
		uci -q set wireless.default_radio1.encryption='sae-mixed'
		uci -q set wireless.default_radio1.ieee80211w='1'
		uci -q set wireless.default_radio1.key="${wifi_key}"
	fi

	uci -q commit wireless
	# 中文：首次启动脚本运行较早，只写入 UCI 配置；后续 network/wifi 启动时会自动应用。
	# 中文：避免在这里主动 reload 无线，防止驱动初始化阶段卡住首次启动流程。
}

setup_default_wifi

setup_adguardhome_dns_port() {
	# 中文：AdGuardHome 不接管 53 端口，避免影响 dnsmasq 的 DNS 与 DHCP。
	# 中文：如用户已经初始化过 AdGuardHome 配置，只把 DNS 监听端口固定到 5335。
	local agh_config

	agh_config="$(uci -q get adguardhome.config.config_file 2>/dev/null || true)"
	[ -n "${agh_config}" ] || agh_config='/etc/adguardhome/adguardhome.yaml'

	if uci -q get adguardhome.config >/dev/null; then
		uci -q set adguardhome.config.enabled='0'
		uci -q set adguardhome.config.config_file="${agh_config}"
		uci -q commit adguardhome
	fi

	# 中文：保留 AdGuardHome 插件和配置入口，但默认不随系统启动，避免占用 53 端口或增加开机负担。
	if [ -x /etc/init.d/adguardhome ]; then
		/etc/init.d/adguardhome disable || true
	fi

	[ -s "${agh_config}" ] || return 0

	awk '
		BEGIN { in_dns = 0; done = 0 }
		/^dns:[[:space:]]*$/ { in_dns = 1; print; next }
		in_dns && /^[^[:space:]][^:]*:/ {
			if (!done) {
				print "  port: 5335"
				done = 1
			}
			in_dns = 0
		}
		in_dns && /^[[:space:]]+port:[[:space:]]*/ {
			if (!done) {
				print "  port: 5335"
				done = 1
			}
			next
		}
		{ print }
		END {
			if (in_dns && !done) {
				print "  port: 5335"
			}
		}
	' "${agh_config}" > "${agh_config}.tmp" && mv "${agh_config}.tmp" "${agh_config}"
}

setup_adguardhome_dns_port

patch_luci_cpuusage_backend() {
	local luci_ucode tmp_file

	luci_ucode='/usr/share/rpcd/ucode/luci'
	[ -s "${luci_ucode}" ] || return 0
	[ -x /sbin/cpuusage ] || return 0
	grep -q '/sbin/cpuusage' "${luci_ucode}" && return 0

	tmp_file="${luci_ucode}.tmp"
	# 中文：OpenWrt 主线默认只返回 CPU 百分比；这里复用 /sbin/cpuusage 输出 CPU 与加速状态。
	awk '
		/const fd = popen\('\''top -n1 \| awk/ {
			print "\t\t\tconst fd = popen('\''/sbin/cpuusage'\'');"
			next
		}
		{ print }
	' "${luci_ucode}" > "${tmp_file}" && mv "${tmp_file}" "${luci_ucode}"
}

patch_luci_cpuusage_backend

setup_optional_services_disabled() {
	# 中文：代理和 DDNS 插件保留配置入口，但默认开关关闭；需要时由用户在 LuCI 中手动启用。
	local ddns_service

	if uci -q get daed.config >/dev/null; then
		uci -q set daed.config.enabled='0'
		uci -q commit daed
	fi

	if uci -q get openclash.config >/dev/null; then
		uci -q set openclash.config.enable='0'
		uci -q commit openclash
	fi

	if uci -q get homeproxy.config >/dev/null; then
		uci -q set homeproxy.config.main_node='nil'
		uci -q set homeproxy.config.main_udp_node='same'
		uci -q commit homeproxy
	fi

	if uci -q get homeproxy.control >/dev/null; then
		uci -q set homeproxy.control.lan_proxy_mode='disabled'
		uci -q commit homeproxy
	fi

	if uci -q get homeproxy.server >/dev/null; then
		uci -q set homeproxy.server.enabled='0'
		uci -q commit homeproxy
	fi

	if uci -q show ddns >/dev/null 2>&1; then
		for ddns_service in $(uci -q show ddns | sed -n "s/^ddns\.\([^.=]*\)=service$/\1/p"); do
			uci -q set "ddns.${ddns_service}.enabled=0" || true
		done
		uci -q commit ddns
	fi

	if [ -x /etc/init.d/ddns ]; then
		/etc/init.d/ddns disable || true
	fi
}

setup_optional_services_disabled

if [ -x /etc/init.d/irqbalance ]; then
	# 中文：irqbalance 的 init.d enable 只控制启动链接；UCI enabled=1 才会真正拉起进程。
	if ! uci -q get irqbalance.irqbalance >/dev/null; then
		uci -q set irqbalance.irqbalance='irqbalance'
	fi
	uci -q set irqbalance.irqbalance.enabled='1'
	uci -q commit irqbalance
	/etc/init.d/irqbalance enable
fi

exit 0
